# .config/sv/redshift/run
#!/bin/sh

exec redshift-gtk -l44.7807:20.5003
